# etc/profile.d/rust-download.sh
export RUSTUP_HOME="/usr/local/lib/rust-download"
export PATH="$RUSTUP_HOME/bin:$PATH"
