# etc/profile.d/rust-download.sh
RUST_DOWNLOAD_DIR="/usr/local/lib/rust-download"
export RUSTUP_HOME="$RUST_DOWNLOAD_DIR/rustup"
SYSTEM_WIDE_CARGO_HOME="$RUST_DOWNLOAD_DIR/cargo"
export PATH="$SYSTEM_WIDE_CARGO_HOME/bin:$PATH"
